# verifuse_v2/deploy/verifuse-vertex.timer
[Unit]
Description=Run VeriFuse V2 Vertex AI processing weekly on Fridays

[Timer]
OnCalendar=Fri *-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
